# contrib/init/bitcoind.conf
description "Gpunion Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GPUNIOND_BIN="/usr/bin/gpuniond"
env GPUNIOND_USER="gpunion"
env GPUNIOND_GROUP="gpunion"
env GPUNIOND_PIDDIR="/var/run/gpuniond"
# upstart can't handle variables constructed with other variables
env GPUNIOND_PIDFILE="/var/run/gpuniond/gpuniond.pid"
env GPUNIOND_CONFIGFILE="/etc/gpunion/gpunion.conf"
env GPUNIOND_DATADIR="/var/lib/gpuniond"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # gpuniond will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GPUNIOND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run gpuniond."
        echo "The setting must appear in $GPUNIOND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Gpunion Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GPUNIOND_PIDDIR"
    chmod 0755 "$GPUNIOND_PIDDIR"
    chown $GPUNIOND_USER:$GPUNIOND_GROUP "$GPUNIOND_PIDDIR"
    chown $GPUNIOND_USER:$GPUNIOND_GROUP "$GPUNIOND_CONFIGFILE"
    chmod 0660 "$GPUNIOND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GPUNIOND_PIDFILE" \
    --chuid $GPUNIOND_USER:$GPUNIOND_GROUP \
    --exec "$GPUNIOND_BIN" \
    -- \
    -pid="$GPUNIOND_PIDFILE" \
    -conf="$GPUNIOND_CONFIGFILE" \
    -datadir="$GPUNIOND_DATADIR" \
    -disablewallet \
    -daemon

